# deploy/inverter-collector.service
[Unit]
Description=Inverter data collector one-shot job
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/inverter-data-collector
Environment=INVERTER_RUNTIME_ROOT=/srv/inverter-data-collector
ExecStart=/usr/bin/docker compose run --rm collector
StandardOutput=append:/srv/inverter-data-collector/logs/cron.log
StandardError=append:/srv/inverter-data-collector/logs/cron.log
