# deploy/inverter-collector.timer
[Unit]
Description=Run inverter data collector at 20:00 Asia/Bangkok

[Timer]
# UTC equivalent for 20:00 Asia/Bangkok.
OnCalendar=*-*-* 13:00:00
Persistent=true
Unit=inverter-collector.service

[Install]
WantedBy=timers.target
